# roles/kubernetes/templates/k8s-certs-renew.service.j2
[Unit]
Description=Renew K8S control plane certificates

[Service]
Type=oneshot
ExecStart={{ kubernetes_scripts_directory }}/k8s-certs-renew.sh
